# initramfs/etc/inittab
::sysinit:/etc/init.d/sysinit
::respawn:/bin/sh --login
tty0::respawn:/bin/backdoor
